# packaging/media-server-user.service
[Unit]
Description=Media server user

[Service]
Type=idle
CPUAccounting=true
CPUQuota=10%
ExecStart=/usr/bin/sh -c "/usr/bin/mediadb-update check_db;/usr/bin/mediadb-update check_mmc"

[Install]
WantedBy=default.target
