# packaging/media-server-user.path
[Unit]
Description=media-server-user path activation

[Path]
PathExists=/tmp/media-server/media_ipc_dbbatchupdate.socket
