# packaging/media-server-user.service
[Unit]
Description=Media server user

[Service]
Type=idle
CPUAccounting=true
CPUQuota=10
ExecStartPre=/usr/bin/sleep 5
ExecStart=/usr/bin/mediadb-update -r %h

[Install]
WantedBy=default.target
