# packaging/media-server.service
[Unit]
Description=Media server

[Service]
Type=simple
ExecStart=/usr/bin/media-server

[Install]
WantedBy=multi-user.target
